# docker/app/etc/crontab/crontab
# Laravel Scheduler
* * * * * cd /var/www && php artisan schedule:run >> /dev/null
